# overlay/etc/supervisor/conf.d/mineos.conf
[program:mineos]
command=/root/.nvm/versions/node/v4.0.0/bin/node webui.js
directory=/usr/games/minecraft
user=root
autostart=true
autorestart=true
redirect_stderr=true
